# dist/bin/jsbeat.service
[Unit]
Description=A light, Javascript based, Beat to collect data for the Open Collector
Documentation=https://github.com/TonyMasse/jsBeat
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=0

[Service]
ExecStart=/opt/jsBeat/bin/start.sh
Restart=always
RestartSec=5
Type=simple

[Install]
WantedBy=multi-user.target
